# scripts/evic-usb-rtc-sync.service
[Unit]
Description=Evic RTC sync

[Service]
ExecStart=/usr/bin/evic-usb time
RemainAfterExit=yes
