# scripts/evic-usb-rtc-sync.service
[Unit]
Description=Evic RTC sync

[Service]
ExecStart=/bin/bash -c 'lsusb -d 0416:5020 && evic-usb time'
RemainAfterExit=yes
